# templates/etc/systemd/system/rpi_exporter.service
[Unit]
Description=RPi Exporter

[Service]
User=nobody
ExecStart=/usr/local/bin/rpi_exporter

[Install]
WantedBy=multi-user.target
